# internet-access-schedule-on.timer
[Unit]
Description=Turn access to dns on at 7:00 am

[Timer]
OnCalendar=7:00:00
Persistent=true

[Install]
WantedBy=timers.target

